# kafka/supervisor/kafka-aio.conf
[program:kafka]
command=/scripts/start-kafka-aio.sh
autostart=true
autorestart=true
