# kafka/supervisor/kdc.conf
[program:producer]
command=/usr/lib/heimdal-servers/kdc --config-file=/etc/heimdal-kdc/kdc.conf -P 88
autostart=true
autorestart=true
