# kafka/supervisor/producer.conf
[program:producer]
command=python3 -u /scripts/producer.py
autostart=true
autorestart=true
