# kafka/supervisor/zookeeper.conf
[program:zookeeper]
command=/scripts/start-zookeeper.sh
autostart=true
autorestart=true
stdout_logfile = /var/log/zookeeper/zookeeper.log
stdout_logfile_backups = 0
redirect_stderr = true
